# infra/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ---------- redis ----------
[program:redis]
command=redis-server --save "" --appendonly no --requirepass %(ENV_REDIS_PASSWORD)s
autorestart=true
stdout_logfile=/var/log/redis-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/redis-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=10

; ---------- postgres ----------
[program:postgres]
command=bash -c "pg_ctlcluster 16 main start --foreground"
autorestart=true
stdout_logfile=/var/log/postgres-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/postgres-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=20
startsecs=5

; ---------- pg-restore ----------
; Runs once on startup: restores from network volume dump, then exits.
[program:pg-restore]
command=/workspace/blotter/infra/postgres/pg_restore.sh
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/var/log/pg-restore.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
stderr_logfile=/var/log/pg-restore-err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1
priority=21

; ---------- ollama ----------
[program:ollama]
command=ollama serve
environment=OLLAMA_HOST="127.0.0.1:11434",OLLAMA_NUM_PARALLEL="2",OLLAMA_KEEP_ALIVE="24h",CUDA_VISIBLE_DEVICES="0"
autorestart=true
stdout_logfile=/var/log/ollama-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/ollama-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=25
startsecs=5

; ---------- cloudflared ----------
[program:cloudflared]
command=cloudflared tunnel run --token %(ENV_CLOUDFLARED_TOKEN)s
autorestart=true
stdout_logfile=/var/log/cloudflared-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/cloudflared-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=30

; ---------- monitoring ----------
[program:monitoring]
command=bash /workspace/blotter/infra/monitoring/monitoring_loop.sh
autorestart=true
stdout_logfile=/var/log/monitoring.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
stderr_logfile=/var/log/monitoring-err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1
priority=5

; ---------- pg-backup ----------
; Hourly pg_dump to network volume. Runs once, sleeps, repeats.
[program:pg-backup]
command=bash -c "while true; do /workspace/blotter/infra/postgres/pg_dump.sh; sleep 3600; done"
autorestart=true
stdout_logfile=/var/log/pg-backup.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
stderr_logfile=/var/log/pg-backup-err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1
priority=22

; ---------- api ----------
[program:api]
command=uv run uvicorn blotter.api:app --host 0.0.0.0 --port 8080
directory=/workspace/blotter/backend
environment=PATH="/root/.local/bin:%(ENV_PATH)s",REDIS_PASSWORD="%(ENV_REDIS_PASSWORD)s"
autorestart=true
stdout_logfile=/var/log/blotter-api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/blotter-api-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=35
startsecs=3

; ---------- pipeline ----------
; Depends on redis + postgres being up. priority=40 starts it last.
[program:pipeline]
command=uv run blotter stream start --openmhz --transcriber-workers 4
directory=/workspace/blotter/backend
environment=PATH="/root/.local/bin:%(ENV_PATH)s",GOOGLE_APPLICATION_CREDENTIALS="/workspace/blotter-gcs-key.json",REDIS_PASSWORD="%(ENV_REDIS_PASSWORD)s",STREAM_AD_SKIP_SECONDS="0",OPENMHZ_SYSTEMS="lapdvalley,lapdwest,chi_cpd,cltp25,philly,psern1,sfp25,pgcomd,pdx2,ntirnd1,nwhc,dane_com,monroecony,dcfd,mnhennco,njicsunion,gcrn,mcbsimcast,sc21102,scpd,snacc,apsp25,bacop25,indydps"
autorestart=true
stdout_logfile=/var/log/blotter-pipeline.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/blotter-pipeline-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=40
startsecs=10
startretries=999
