# infra/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ---------- redis ----------
[program:redis]
command=redis-server --save "" --appendonly no --requirepass %(ENV_REDIS_PASSWORD)s
autorestart=true
stdout_logfile=/var/log/redis-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/redis-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=10

; ---------- clickhouse ----------
; --no-daemon keeps it in foreground so supervisord can manage it
[program:clickhouse]
command=clickhouse-server --config-file=/etc/clickhouse-server/config.xml
autorestart=true
stdout_logfile=/var/log/clickhouse-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
stderr_logfile=/var/log/clickhouse-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=20
startsecs=5

; ---------- ollama ----------
[program:ollama]
command=ollama serve
environment=OLLAMA_HOST="127.0.0.1:11434",OLLAMA_NUM_PARALLEL="2",OLLAMA_KEEP_ALIVE="24h"
autorestart=true
stdout_logfile=/var/log/ollama-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/ollama-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=25
startsecs=5

; ---------- cloudflared ----------
[program:cloudflared]
command=cloudflared tunnel --config /workspace/blotter/infra/cloudflared/config.yml run
autorestart=true
stdout_logfile=/var/log/cloudflared-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/cloudflared-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=30

; ---------- pipeline ----------
; Depends on redis + clickhouse being up. priority=40 starts it last.
[program:pipeline]
command=uv run blotter stream start --openmhz --transcriber-workers 5
directory=/workspace/blotter/backend
environment=PATH="/root/.local/bin:%(ENV_PATH)s",GOOGLE_APPLICATION_CREDENTIALS="/workspace/blotter-gcs-key.json",REDIS_PASSWORD="%(ENV_REDIS_PASSWORD)s",STREAM_AD_SKIP_SECONDS="0",OPENMHZ_SYSTEMS="lapdvalley,lapdwest,chi_cpd,cltp25,philly,psern1,sfp25,pgcomd,pdx2,ntirnd1,nwhc,dane_com,monroecony,dcfd,mnhennco,njicsunion,gcrn,mcbsimcast,sc21102,scpd,snacc,apsp25,bacop25,indydps"
autorestart=true
stdout_logfile=/var/log/blotter-pipeline.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/blotter-pipeline-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=40
startsecs=10
startretries=999
